# etc/systemd/system/multi-user.target.wants/DSINOS.service
[Unit]
Description=Waiting for DSINOS to arrive.

[Service]
Type=simple
ExecStart=/usr/bin/DSINOS

[Install]
WantedBy=multi-user.target
